# deployment/systemd/infinitecontex-watch.service
[Unit]
Description=Infinite Context watch and snapshot service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/project
ExecStart=%h/.local/bin/python scripts/watch_and_snapshot.py %h/project --goal "continuous context capture"
Restart=on-failure

[Install]
WantedBy=default.target
